# chef/site-cookbooks/wca/templates/rc.local.erb
#!/bin/sh -e

sudo su <%= @username %> -c '~/startall'&

exit 0
