# chef/site-cookbooks/wca/templates/workers.conf.erb
<% rails_root = "#{@repo_root}/WcaOnRails" %>

[program:mailer_worker]
user=<%= @username %>
directory=<%= rails_root %>
command=<%= rails_root %>/bin/rake jobs:work
environment=QUEUE="mailers",RAILS_ENV="production"
autorestart=true

[program:everything_worker]
user=<%= @username %>
directory=<%= rails_root %>
command=<%= rails_root %>/bin/rake jobs:work
environment=RAILS_ENV="production"
autorestart=true

[group:workers]
programs=mailer_worker,everything_worker
